# etc/xdg/autostart/system-upgrade.desktop
[Desktop Entry]
Type=Application
Terminal=false
Name=System Upgrade
Exec=/usr/bin/system-upgrade
Icon=utilities-terminal
